# install.d/inittab

## BEGIN-NETPASS
# Run npsvc in runlevel 3
npsv:3:respawn:/opt/netpass/bin/npsvc.pl -m %SMTPHOST%
## END-NETPASS
